# ansible/supervisor_ksql.conf
[program:ksql]
command=/opt/ksql/bin/ksql-server-start /opt/ksql/config/ksqlserver.properties
user=root